# src/etc/service/nginx/run
#!/bin/sh
set -e
[ -f /env ] && . /env 2>/dev/null
exec /usr/sbin/nginx -g 'daemon off;' >/dev/null
