# src/etc/service/php-fpm/run
#!/bin/sh
set -e
[ -f /env ] && . /env 2>/dev/null
exec php-fpm --nodaemonize
