# src/etc/service/php-fpm/run
#!/bin/sh

source /env

php-fpm --nodaemonize &
pid=$!

trap 'echo "php-fpm shutting down..."; kill -TERM "$pid"; wait "$pid"; exit 0' TERM INT

wait "$pid"
